# infra/supervisord.conf
; ═══════════════════════════════════════════════════════════════
;  SENTINEL V5.5 — Sovereign Process Supervisor
; ═══════════════════════════════════════════════════════════════
;
;  Lifecycle Dependency:
;    sentinel-sidecar starts first (startsecs=5).
;    python-worker starts second (startsecs=10) — passive gate
;    ensures the UDS is available before the first arbitration call.
;
;  2 AM Recovery:
;    autorestart=true on both programs. If the sidecar crashes,
;    supervisord restarts it. The Python worker's get_fast_veredict
;    must implement retry-with-backoff on ConnectionRefusedError
;    to survive the ~3s sidecar reboot window.
;
;  Deploy:
;    supervisord -c /etc/sentinel/supervisord.conf
;
; ═══════════════════════════════════════════════════════════════

[unix_http_server]
file=/tmp/supervisor.sock   ; path to the socket file

[supervisord]
logfile=/var/log/supervisord.log ; main log file
pidfile=/var/run/supervisord.pid ; supervisord pidfile
nodaemon=false               ; start in foreground if true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL for a unix socket

; ─────────────────────────────────────────────────────
;  SENTINEL SIDECAR — Go UDS Daemon
;
;  - startsecs=5:       Must survive 5s before considered "running"
;  - stopwaitsecs=10:   Allows final snapshot + WAL flush on SIGTERM
;  - autorestart=true:  Closes the "2 AM crash" gap
; ─────────────────────────────────────────────────────

[program:sentinel-sidecar]
command=/usr/local/bin/sentinel-sidecar
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
user=sentinel-admin
stderr_logfile=/var/log/sentinel-sidecar.err.log
stdout_logfile=/var/log/sentinel-sidecar.out.log
; Hub URL for background sync (Blocker #4) + WAL path for replay (Blocker #3)
environment=SENTINEL_HUB_URL="https://hub.sentinel.internal",SENTINEL_WAL_PATH="/var/sentinel/wal.jsonl",SENTINEL_GRAPH_SNAPSHOT="/var/sentinel/graph.json",SENTINEL_SOCKET_GROUP="sentinel-runners",SENTINEL_SIDECAR_SOCKET="/tmp/sentinel_veritas.sock"

; ─────────────────────────────────────────────────────
;  PYTHON WORKER — Governed Asset Processor
;
;  - startsecs=10:      Passive dependency gate — ensures the
;                       sidecar has had 5+ seconds to bind the UDS
;                       before the worker's first arbitration call.
;  - autorestart=true:  Worker restarts on crash.
; ─────────────────────────────────────────────────────

[program:python-worker]
command=python /opt/sentinel/remove_bg_governed.py
autostart=true
autorestart=true
startsecs=10
user=sentinel-runner
environment=SENTINEL_TENANT_ID="TMU_VOSPITAL_001",SENTINEL_SOVEREIGN_SCOPE="TMU_VOSPITAL",SENTINEL_SIDECAR_SOCKET="/tmp/sentinel_veritas.sock"
stderr_logfile=/var/log/python-worker.err.log
stdout_logfile=/var/log/python-worker.out.log
